# bot.service
[Unit]
Description=Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /path/to/your/bot.py
Restart=on-failure

[Install]
WantedBy=default.target
